# app/supervisor.d/uwsgi.conf
[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /var/www/uwsgi.ini
autostart=true
autorestart=true
stdout_logfile=/var/www/logs/uwsgi_out.log
stderr_logfile=/var/www/logs/uwsgi_err.log